# gitbuddy.service
[Unit]
Description=GitBuddy Repository Sync Service
Documentation=https://github.com/yourusername/git-buddy
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.config/git-buddy/git_puller_service.py
WorkingDirectory=%h/.config/git-buddy/
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target